# vodarchiver.service
# /etc/systemd/system/vodarchiver.service
[Unit]
Description=Vod Archiver
After=network.target

[Service]
Environment=PATH=<path to bun bin>:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Type=simple
WorkingDirectory=<app path>
ExecStart=<path to bun> run index.ts
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target